# indexing/crontab
0 1 * * 0,2-6 INDEX_DELETE=false /usr/src/app/cron.sh
0 1 * * 1 INDEX_DELETE=true /usr/src/app/cron.sh
#